# distributions/docker/compose/vnfm/generic/supervisord.conf
[supervisord]
nodaemon=true

[program:update-vnfm-generic-config]
command=sh -c "export RABBITMQ_BROKERIP=${RABBITMQ_BROKERIP} && sed -i -e \"s/vnfm.rabbitmq.brokerIp\s*=\s*localhost/vnfm.rabbitmq.brokerIp = ${RABBITMQ_BROKERIP}/g\" /etc/openbaton/openbaton-vnfm-generic.properties"
autorestart=false
startsecs=0

[program:openbaton-vnfm-generic]
command=/usr/bin/java -jar /usr/lib/openbaton/vnfm/generic/openbaton-vnfm-generic-3.2.0.jar --spring.config.location=file:/etc/openbaton/openbaton-vnfm-generic.properties
