# distributions/docker/compose/nfvo/supervisord.conf
[supervisord]
nodaemon=true

; Only if MySQL has been installed during the bootstrap 
;[program:mysql-server]
;command=/etc/init.d/mysql start

[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server

[program:update-nfvo-config]
command=sh -c "export RABBITMQ_BROKERIP=${RABBITMQ_BROKERIP} && sed -i -e \"s/nfvo.rabbit.brokerIp\s*=\s*localhost/nfvo.rabbit.brokerIp = ${RABBITMQ_BROKERIP}/g\" /etc/openbaton/openbaton-nfvo.properties"
autorestart=false
startsecs=0

[program:openbaton-nfvo]
command=/usr/bin/java -jar /usr/lib/openbaton/nfvo/openbaton-nfvo-3.2.0.jar --spring.config.location=file:/etc/openbaton/openbaton-nfvo.properties --no-console
