# distributions/docker/develop/supervisord.conf
[supervisord]
nodaemon=true

[program:update-nfvo-config]
command=sh -c "export RABBITMQ_BROKERIP=${RABBITMQ_BROKERIP} && sed -i -e \"s/nfvo.rabbit.brokerIp\s*=\s*localhost/nfvo.rabbit.brokerIp = ${RABBITMQ_BROKERIP}/g\" /etc/openbaton/openbaton-nfvo.properties"
autorestart=false
startsecs=0

[program:update-vnfm-generic-config]
command=sh -c "export RABBITMQ_BROKERIP=${RABBITMQ_BROKERIP} && sed -i -e \"s/vnfm.rabbitmq.brokerIp\s*=\s*localhost/vnfm.rabbitmq.brokerIp = ${RABBITMQ_BROKERIP}/g\" /etc/openbaton/vnfm/generic/openbaton-vnfm-generic.properties"
autorestart=false
startsecs=0

[program:mysql-server]
command=service mysql start

[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server

[program:openbaton-nfvo]
command=/usr/bin/java -jar /usr/lib/openbaton/nfvo/openbaton-nfvo-3.2.0.jar --spring.config.location=file:/etc/openbaton/openbaton-nfvo.properties --no-console

[program:openbaton-vnfm-generic]
command=/usr/bin/java -jar /usr/lib/openbaton/vnfm/generic/openbaton-vnfm-generic-3.2.0.jar --spring.config.location=file:/etc/openbaton/openbaton-vnfm-generic.properties
